# access-point/looped-keypad.service
[Unit]
Description=Looped Keypad DTMF Scanner
After=network.target
Wants=looped-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/admin/Looped/frontend/backend/app
ExecStart=/usr/bin/python3 /home/admin/Looped/frontend/backend/app/keypad_dtmf.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
